# cowrie/supervisor/cowrie.conf
[program:cowrie]
command=/cowrie/cowrie-git/bin/cowrie start
directory=/cowrie/cowrie-git
user=cowrie
autorestart=true
redirect_stderr=true
